# docker/supervisor/golang.ini
[program:golang]
process_name=%(program_name)s
directory=/go/src/api
command=bee run -downdoc=true -gendoc=true
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/stdout_golang.log
stderr_logfile=/var/log/supervisor/stderr_golang.log
